# resources/supervisor/pod.conf
[program:{{tag}}]
command=python /opt/{{tag}}/pod/pod.py