# resources/supervisor/pod.conf
[program:crow]
command=python /opt/crow/pod/pod.py